# rootfs/etc/services.d/ecoflow-panel/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# EcoFlow Panel — s6 service runner.
# Translates the add-on's Options (Configuration tab in HA) into the env vars
# the Node server expects, then execs npm start.
# ==============================================================================
set -e

# ── Required EcoFlow API credentials ─────────────────────────────────────────
export ECOFLOW_ACCESS_KEY="$(bashio::config 'ECOFLOW_ACCESS_KEY')"
export ECOFLOW_SECRET_KEY="$(bashio::config 'ECOFLOW_SECRET_KEY')"
export ECOFLOW_API_HOST="$(bashio::config 'ECOFLOW_API_HOST')"

if bashio::var.is_empty "${ECOFLOW_ACCESS_KEY}" || bashio::var.is_empty "${ECOFLOW_SECRET_KEY}"; then
  bashio::log.fatal "ECOFLOW_ACCESS_KEY and ECOFLOW_SECRET_KEY must be set in"
  bashio::log.fatal "the add-on Configuration tab. Set them and restart."
  # sleep then exit so HA doesn't restart-loop us into the log
  sleep 30
  exit 1
fi

# ── Server + persistence ─────────────────────────────────────────────────────
export PORT=8787
export HOST=::
export DB_PATH=/data/ecoflow.db
export WEB_DIST_PATH=/app/web/dist
export LOG_LEVEL="$(bashio::config 'LOG_LEVEL')"
export FORECAST_LAT="$(bashio::config 'FORECAST_LAT')"
export FORECAST_LON="$(bashio::config 'FORECAST_LON')"

# ── Telnet "control-room" TUI ────────────────────────────────────────────────
if bashio::config.true 'TELNET_ENABLED'; then
  export TELNET_ENABLED=1
else
  export TELNET_ENABLED=0
fi
export TELNET_HOST=::
export TELNET_PORT="$(bashio::config 'TELNET_PORT')"

# ── Push notifications ───────────────────────────────────────────────────────
export NOTIFY_CHANNEL="$(bashio::config 'NOTIFY_CHANNEL')"
export NOTIFY_MIN_SEVERITY="$(bashio::config 'NOTIFY_MIN_SEVERITY')"
if bashio::config.true 'NOTIFY_RESOLVED'; then
  export NOTIFY_RESOLVED=1
else
  export NOTIFY_RESOLVED=0
fi
export NOTIFY_NTFY_SERVER="$(bashio::config 'NOTIFY_NTFY_SERVER')"
export NOTIFY_NTFY_TOPIC="$(bashio::config 'NOTIFY_NTFY_TOPIC')"
export NOTIFY_PUSHOVER_TOKEN="$(bashio::config 'NOTIFY_PUSHOVER_TOKEN')"
export NOTIFY_PUSHOVER_USER="$(bashio::config 'NOTIFY_PUSHOVER_USER')"
export NOTIFY_WEBHOOK_URL="$(bashio::config 'NOTIFY_WEBHOOK_URL')"

# ── v0.7.5 quiet-hours digest ────────────────────────────────────────────────
export NOTIFY_QUIET_HOURS="$(bashio::config 'NOTIFY_QUIET_HOURS')"
export NOTIFY_DIGEST_HOUR="$(bashio::config 'NOTIFY_DIGEST_HOUR')"

# ── v0.7.5 NWS storm-preparedness (opt-in, US-only) ──────────────────────────
if bashio::config.true 'NWS_ENABLED'; then
  export NWS_ENABLED=1
else
  export NWS_ENABLED=0
fi

# ── v0.7.5 HA MQTT Discovery (opt-in) ────────────────────────────────────────
if bashio::config.true 'MQTT_DISCOVERY_ENABLED'; then
  export MQTT_DISCOVERY_ENABLED=1
else
  export MQTT_DISCOVERY_ENABLED=0
fi
export MQTT_DISCOVERY_HOST="$(bashio::config 'MQTT_DISCOVERY_HOST')"
export MQTT_DISCOVERY_PORT="$(bashio::config 'MQTT_DISCOVERY_PORT')"
export MQTT_DISCOVERY_USER="$(bashio::config 'MQTT_DISCOVERY_USER')"
export MQTT_DISCOVERY_PASS="$(bashio::config 'MQTT_DISCOVERY_PASS')"
export MQTT_DISCOVERY_PREFIX="$(bashio::config 'MQTT_DISCOVERY_PREFIX')"

# v0.9.9 — diagnostic write-command endpoint. Off unless this is non-empty.
# Used to probe undocumented EcoFlow commands (e.g. SHP2 reboot, which
# isn't in the public IoT Open API).
export WRITE_DEBUG_TOKEN="$(bashio::config 'WRITE_DEBUG_TOKEN')"
if bashio::config.has_value 'WRITE_DEBUG_TOKEN'; then
  bashio::log.warning "WRITE_DEBUG_TOKEN is SET — POST /api/device/send-command is enabled."
fi

# v0.9.18 — ship-wide audible broadcast to HomePod / Sonos.
if bashio::config.true 'BROADCAST_ENABLED'; then
  export BROADCAST_ENABLED=true
else
  export BROADCAST_ENABLED=false
fi
export BROADCAST_TARGETS="$(bashio::config 'BROADCAST_TARGETS')"
# v1.25.0 — SIP / announce-only targets (e.g. the Switchboard cordless) that get
# alarms via media_player.play_media(announce) instead of Music Assistant.
export BROADCAST_SIP_TARGETS="$(bashio::config 'BROADCAST_SIP_TARGETS')"
export BROADCAST_AUDIO_BASE="$(bashio::config 'BROADCAST_AUDIO_BASE')"
export BROADCAST_VOLUME="$(bashio::config 'BROADCAST_VOLUME')"
# v0.15.4 — announcement reliability/repeat knobs (ecobee-friendly)
export BROADCAST_REPEAT="$(bashio::config 'BROADCAST_REPEAT')"
export BROADCAST_REPEAT_GAP_MS="$(bashio::config 'BROADCAST_REPEAT_GAP_MS')"
export BROADCAST_ANNOUNCE_VOLUME="$(bashio::config 'BROADCAST_ANNOUNCE_VOLUME')"
export BROADCAST_ANNOUNCE_RETRIES="$(bashio::config 'BROADCAST_ANNOUNCE_RETRIES')"
if bashio::config.true 'BROADCAST_USE_PRE_ANNOUNCE'; then
  export BROADCAST_USE_PRE_ANNOUNCE=true
else
  export BROADCAST_USE_PRE_ANNOUNCE=false
fi
export BROADCAST_MIN_SEVERITY="$(bashio::config 'BROADCAST_MIN_SEVERITY')"
export BROADCAST_QUIET_HOURS="$(bashio::config 'BROADCAST_QUIET_HOURS')"
# v0.23.0 — opt-in: critical alerts break through quiet hours (push + audible).
# Default false ⇒ quiet hours silence every tier (alert still shows on-screen +
# morning digest). Shared by alertMonitor (push) and broadcast (audible).
if bashio::config.true 'CRITICAL_BREAKS_QUIET_HOURS'; then
  export CRITICAL_BREAKS_QUIET_HOURS=true
else
  export CRITICAL_BREAKS_QUIET_HOURS=false
fi
# v0.23.0 — these were declared in config.yaml options/schema but never exported
# here, so the HA Configuration-tab knobs were DEAD (the server only ever saw its
# code defaults). Lead-silence especially matters: Music Assistant 2.9's faster
# AirPlay RAOP start needs a longer lead-in to avoid clipping the chime's start.
# (Removed the stale BROADCAST_TTS_SERVICE/TTS_LANGUAGE/SONOS_RESTORE/
# USE_MUSIC_ASSISTANT exports — those options no longer exist in config.yaml.)
export BROADCAST_LEAD_SILENCE_MS="$(bashio::config 'BROADCAST_LEAD_SILENCE_MS')"
export BROADCAST_CHIME_GAP_MS="$(bashio::config 'BROADCAST_CHIME_GAP_MS')"
export BROADCAST_CHIME_PACK="$(bashio::config 'BROADCAST_CHIME_PACK')"
export BROADCAST_WYOMING_HOST="$(bashio::config 'BROADCAST_WYOMING_HOST')"
export BROADCAST_WYOMING_PORT="$(bashio::config 'BROADCAST_WYOMING_PORT')"
export BROADCAST_WYOMING_VOICE="$(bashio::config 'BROADCAST_WYOMING_VOICE')"
# v0.62.0 — bilingual second pass (English then Spanish).
if bashio::config.true 'BROADCAST_BILINGUAL'; then
  export BROADCAST_BILINGUAL=true
else
  export BROADCAST_BILINGUAL=false
fi
export BROADCAST_WYOMING_VOICE_ES="$(bashio::config 'BROADCAST_WYOMING_VOICE_ES')"
if bashio::config.true 'BROADCAST_ENABLED'; then
  bashio::log.info "broadcast: enabled (targets: $BROADCAST_TARGETS)"
fi
# v0.15.2 — off-grid honesty + intelligent load-shedding ADVISOR (read + advise;
# HA automations actuate). GRID_AVAILABLE must be false for an islanded site.
if bashio::config.true 'GRID_AVAILABLE'; then
  export GRID_AVAILABLE=true
else
  export GRID_AVAILABLE=false
fi
# v0.23.0 — optional HA entity reporting real grid presence (grid-aware floor).
export GRID_PRESENCE_ENTITY="$(bashio::config 'GRID_PRESENCE_ENTITY')"
# v1.6.0 — host power self-monitor: an HA binary_sensor (device_class=problem)
# reporting the Pi's kernel under-voltage flag. When set, the alarm engine warns
# on a browning-out host before it takes the whole monitor dark. Empty = dormant.
export HOST_POWER_ENTITY="$(bashio::config 'HOST_POWER_ENTITY')"
# v0.72.0 — optional SN->HA-ping-entity JSON map (cloud-wedge vs real-outage detection).
export ECOFLOW_DEVICE_REACHABILITY="$(bashio::config 'ECOFLOW_DEVICE_REACHABILITY')"
# v0.23.0 — wire the battery-alarm enable toggles + re-announce cadence to env.
# These were declared in config.yaml but never exported here, so the HA toggles
# were dead and the server only ever saw its code defaults (both enabled). Now
# the Configuration-tab switches actually take effect.
if bashio::config.true 'BATTERY_SOC_ALARM_ENABLED'; then
  export BATTERY_SOC_ALARM_ENABLED=true
else
  export BATTERY_SOC_ALARM_ENABLED=false
fi
if bashio::config.true 'BATTERY_RUNWAY_ALARM_ENABLED'; then
  export BATTERY_RUNWAY_ALARM_ENABLED=true
else
  export BATTERY_RUNWAY_ALARM_ENABLED=false
fi
export BATTERY_RUNWAY_ALARM_REANNOUNCE_MIN="$(bashio::config 'BATTERY_RUNWAY_ALARM_REANNOUNCE_MIN')"
if bashio::config.true 'LOAD_SHEDDING_ADVISORY_ENABLED'; then
  export LOAD_SHEDDING_ADVISORY_ENABLED=true
else
  export LOAD_SHEDDING_ADVISORY_ENABLED=false
fi
export LOAD_SHEDDING_SHED_ENTITIES="$(bashio::config 'LOAD_SHEDDING_SHED_ENTITIES')"
export LOAD_SHEDDING_RUNWAY_THRESHOLD_H="$(bashio::config 'LOAD_SHEDDING_RUNWAY_THRESHOLD_H')"
export LOAD_SHEDDING_RESTORE_MARGIN_H="$(bashio::config 'LOAD_SHEDDING_RESTORE_MARGIN_H')"
# SUPERVISOR_TOKEN is provided automatically by the HA Supervisor —
# no need to export it manually. The Node process inherits it.

bashio::log.info "Starting EcoFlow Panel on :${PORT} (telnet :${TELNET_PORT}, db ${DB_PATH})…"

cd /app/server
exec npm start
