# rootfs/etc/services.d/ecoflow-panel/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# EcoFlow Panel — s6 service runner.
# Translates the add-on's Options (Configuration tab in HA) into the env vars
# the Node server expects, then execs npm start.
# ==============================================================================
set -e

# ── Required EcoFlow API credentials ─────────────────────────────────────────
export ECOFLOW_ACCESS_KEY="$(bashio::config 'ECOFLOW_ACCESS_KEY')"
export ECOFLOW_SECRET_KEY="$(bashio::config 'ECOFLOW_SECRET_KEY')"
export ECOFLOW_API_HOST="$(bashio::config 'ECOFLOW_API_HOST')"

if bashio::var.is_empty "${ECOFLOW_ACCESS_KEY}" || bashio::var.is_empty "${ECOFLOW_SECRET_KEY}"; then
  bashio::log.fatal "ECOFLOW_ACCESS_KEY and ECOFLOW_SECRET_KEY must be set in"
  bashio::log.fatal "the add-on Configuration tab. Set them and restart."
  # sleep then exit so HA doesn't restart-loop us into the log
  sleep 30
  exit 1
fi

# ── Server + persistence ─────────────────────────────────────────────────────
export PORT=8787
export HOST=::
export DB_PATH=/data/ecoflow.db
export WEB_DIST_PATH=/app/web/dist
export LOG_LEVEL="$(bashio::config 'LOG_LEVEL')"
export FORECAST_LAT="$(bashio::config 'FORECAST_LAT')"
export FORECAST_LON="$(bashio::config 'FORECAST_LON')"

# ── Telnet "control-room" TUI ────────────────────────────────────────────────
if bashio::config.true 'TELNET_ENABLED'; then
  export TELNET_ENABLED=1
else
  export TELNET_ENABLED=0
fi
export TELNET_HOST=::
export TELNET_PORT="$(bashio::config 'TELNET_PORT')"

# ── Push notifications ───────────────────────────────────────────────────────
export NOTIFY_CHANNEL="$(bashio::config 'NOTIFY_CHANNEL')"
export NOTIFY_MIN_SEVERITY="$(bashio::config 'NOTIFY_MIN_SEVERITY')"
if bashio::config.true 'NOTIFY_RESOLVED'; then
  export NOTIFY_RESOLVED=1
else
  export NOTIFY_RESOLVED=0
fi
export NOTIFY_NTFY_SERVER="$(bashio::config 'NOTIFY_NTFY_SERVER')"
export NOTIFY_NTFY_TOPIC="$(bashio::config 'NOTIFY_NTFY_TOPIC')"
export NOTIFY_PUSHOVER_TOKEN="$(bashio::config 'NOTIFY_PUSHOVER_TOKEN')"
export NOTIFY_PUSHOVER_USER="$(bashio::config 'NOTIFY_PUSHOVER_USER')"
export NOTIFY_WEBHOOK_URL="$(bashio::config 'NOTIFY_WEBHOOK_URL')"

# ── v0.7.5 quiet-hours digest ────────────────────────────────────────────────
export NOTIFY_QUIET_HOURS="$(bashio::config 'NOTIFY_QUIET_HOURS')"
export NOTIFY_DIGEST_HOUR="$(bashio::config 'NOTIFY_DIGEST_HOUR')"

# ── v0.7.5 NWS storm-preparedness (opt-in, US-only) ──────────────────────────
if bashio::config.true 'NWS_ENABLED'; then
  export NWS_ENABLED=1
else
  export NWS_ENABLED=0
fi

# ── v0.7.5 HA MQTT Discovery (opt-in) ────────────────────────────────────────
if bashio::config.true 'MQTT_DISCOVERY_ENABLED'; then
  export MQTT_DISCOVERY_ENABLED=1
else
  export MQTT_DISCOVERY_ENABLED=0
fi
export MQTT_DISCOVERY_HOST="$(bashio::config 'MQTT_DISCOVERY_HOST')"
export MQTT_DISCOVERY_PORT="$(bashio::config 'MQTT_DISCOVERY_PORT')"
export MQTT_DISCOVERY_USER="$(bashio::config 'MQTT_DISCOVERY_USER')"
export MQTT_DISCOVERY_PASS="$(bashio::config 'MQTT_DISCOVERY_PASS')"
export MQTT_DISCOVERY_PREFIX="$(bashio::config 'MQTT_DISCOVERY_PREFIX')"

# v0.9.9 — diagnostic write-command endpoint. Off unless this is non-empty.
# Used to probe undocumented EcoFlow commands (e.g. SHP2 reboot, which
# isn't in the public IoT Open API).
export WRITE_DEBUG_TOKEN="$(bashio::config 'WRITE_DEBUG_TOKEN')"
if bashio::config.has_value 'WRITE_DEBUG_TOKEN'; then
  bashio::log.warning "WRITE_DEBUG_TOKEN is SET — POST /api/device/send-command is enabled."
fi

# v0.9.18 — ship-wide audible broadcast to HomePod / Sonos.
if bashio::config.true 'BROADCAST_ENABLED'; then
  export BROADCAST_ENABLED=true
else
  export BROADCAST_ENABLED=false
fi
export BROADCAST_TARGETS="$(bashio::config 'BROADCAST_TARGETS')"
export BROADCAST_AUDIO_BASE="$(bashio::config 'BROADCAST_AUDIO_BASE')"
export BROADCAST_VOLUME="$(bashio::config 'BROADCAST_VOLUME')"
# v0.15.4 — announcement reliability/repeat knobs (ecobee-friendly)
export BROADCAST_REPEAT="$(bashio::config 'BROADCAST_REPEAT')"
export BROADCAST_REPEAT_GAP_MS="$(bashio::config 'BROADCAST_REPEAT_GAP_MS')"
export BROADCAST_ANNOUNCE_VOLUME="$(bashio::config 'BROADCAST_ANNOUNCE_VOLUME')"
export BROADCAST_ANNOUNCE_RETRIES="$(bashio::config 'BROADCAST_ANNOUNCE_RETRIES')"
if bashio::config.true 'BROADCAST_USE_PRE_ANNOUNCE'; then
  export BROADCAST_USE_PRE_ANNOUNCE=true
else
  export BROADCAST_USE_PRE_ANNOUNCE=false
fi
export BROADCAST_MIN_SEVERITY="$(bashio::config 'BROADCAST_MIN_SEVERITY')"
export BROADCAST_QUIET_HOURS="$(bashio::config 'BROADCAST_QUIET_HOURS')"
export BROADCAST_TTS_SERVICE="$(bashio::config 'BROADCAST_TTS_SERVICE')"
export BROADCAST_TTS_LANGUAGE="$(bashio::config 'BROADCAST_TTS_LANGUAGE')"
if bashio::config.true 'BROADCAST_SONOS_RESTORE'; then
  export BROADCAST_SONOS_RESTORE=true
else
  export BROADCAST_SONOS_RESTORE=false
fi
# v0.9.23 — MA backend selection
export BROADCAST_USE_MUSIC_ASSISTANT="$(bashio::config 'BROADCAST_USE_MUSIC_ASSISTANT')"
if bashio::config.true 'BROADCAST_ENABLED'; then
  bashio::log.info "broadcast: enabled (targets: $BROADCAST_TARGETS)"
fi
# v0.15.2 — off-grid honesty + intelligent load-shedding ADVISOR (read + advise;
# HA automations actuate). GRID_AVAILABLE must be false for an islanded site.
if bashio::config.true 'GRID_AVAILABLE'; then
  export GRID_AVAILABLE=true
else
  export GRID_AVAILABLE=false
fi
if bashio::config.true 'LOAD_SHEDDING_ADVISORY_ENABLED'; then
  export LOAD_SHEDDING_ADVISORY_ENABLED=true
else
  export LOAD_SHEDDING_ADVISORY_ENABLED=false
fi
export LOAD_SHEDDING_SHED_ENTITIES="$(bashio::config 'LOAD_SHEDDING_SHED_ENTITIES')"
export LOAD_SHEDDING_RUNWAY_THRESHOLD_H="$(bashio::config 'LOAD_SHEDDING_RUNWAY_THRESHOLD_H')"
export LOAD_SHEDDING_RESTORE_MARGIN_H="$(bashio::config 'LOAD_SHEDDING_RESTORE_MARGIN_H')"
# SUPERVISOR_TOKEN is provided automatically by the HA Supervisor —
# no need to export it manually. The Node process inherits it.

bashio::log.info "Starting EcoFlow Panel on :${PORT} (telnet :${TELNET_PORT}, db ${DB_PATH})…"

cd /app/server
exec npm start
