# launchd/com.local.ecoflow-panel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Optional macOS LaunchAgent to auto-start the backend on login.
  Replace /ABSOLUTE/PATH/TO with the absolute path to your clone of this repo
  (run `pwd` from the repo root), then see README → "Auto-start on Mac login".
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.ecoflow-panel</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/npm</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/TO/ecoflow-panel/server</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/ABSOLUTE/PATH/TO/ecoflow-panel/data/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/ABSOLUTE/PATH/TO/ecoflow-panel/data/launchd.err.log</string>
</dict>
</plist>
